# supervisord.conf
[unix_http_server]
file=/run/supervisord.sock   ; (the path to the socket file)

[supervisord]
logfile=/var/log/supervisord.log ; (main log file;default $CWD/supervisord.log)
loglevel=info                ; (log level;default info; others: debug,warn,trace)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock ; use a unix:// URL  for a unix socket


[program:gunicorn]
command = gunicorn --bind 0.0.0.0:8000 project.wsgi --timeout 90 --log-level=DEBUG
directory = /usr/src/app/
stdout_logfile = /usr/src/app/gunicorn_supervisor.log
stderr_logfile = /usr/src/app/gunicorn_supervisor.log

[include]
files = /etc/supervisor.d/*.ini